# vms-lite/packaging/postinst
#!/bin/sh
set -e

# --- LilinVision vms-lite Post-Installation Script ---

USER_NAME="vms-lite"
GROUP_NAME="vms-lite"
DATA_DIR="/var/lib/vms-lite"
VMS_LIBDIR="/usr/lib/vms-lite"
VMS_CONF="/etc/ld.so.conf.d/vms-lite.conf"

case "$1" in
    configure)
        # 1. Create vms-lite group if it doesn't exist
        if ! getent group "$GROUP_NAME" >/dev/null; then
            echo "Creating group $GROUP_NAME..."
            groupadd --system "$GROUP_NAME"
        fi

        # 2. Create vms-lite user if it doesn't exist
        if ! getent passwd "$USER_NAME" >/dev/null; then
            echo "Creating system user $USER_NAME..."
            useradd --system \
                --gid "$GROUP_NAME" \
                --home-dir "$DATA_DIR" \
                --no-create-home \
                --shell /usr/sbin/nologin \
                "$USER_NAME"
        fi

        # 3. Create and set permissions for the data directory
        if [ ! -d "$DATA_DIR" ]; then
            echo "Creating data directory $DATA_DIR..."
            mkdir -p "$DATA_DIR"
        fi
        
        echo "Setting permissions for $DATA_DIR..."
        chown -R "$USER_NAME":"$GROUP_NAME" "$DATA_DIR"
        chmod 750 "$DATA_DIR"

        # 4. Add user to video group to access cameras if needed
        if getent group video >/dev/null; then
            echo "Adding $USER_NAME to video group..."
            usermod -a -G video "$USER_NAME"
        fi

        # 5. Add vms-lite to the real user's group so it can browse their home directory
        REAL_USER="${SUDO_USER:-$(whoami)}"
        if [ "$REAL_USER" != "root" ] && getent passwd "$REAL_USER" >/dev/null; then
            REAL_GROUP=$(id -gn "$REAL_USER")
            echo "Adding $USER_NAME to group $REAL_GROUP (for user $REAL_USER)..."
            usermod -a -G "$REAL_GROUP" "$USER_NAME"
        fi

        # 6. Register private library directory
        echo "Registering private library directory: $VMS_LIBDIR"
        echo "$VMS_LIBDIR" > "$VMS_CONF"

        # Refresh shared library cache
        echo "Updating shared library cache..."
        ldconfig

        # 6. System & Hardware Diagnostics
        echo "[vms-lite] Gathering system information..."
        
        # Virtualization/Environment
        VIRT="Bare Metal / Physical Machine"
        if [ -f /proc/version ] && grep -qi "microsoft" /proc/version; then
            VIRT="WSL (Windows Subsystem for Linux)"
        elif command -v systemd-detect-virt >/dev/null 2>&1; then
            VIRT_TYPE=$(systemd-detect-virt || echo "none")
            if [ "$VIRT_TYPE" != "none" ]; then
                VIRT="Virtual Machine ($VIRT_TYPE)"
            fi
        fi
        echo "  Environment: $VIRT"
        
        # CPU Info
        if command -v lscpu >/dev/null 2>&1; then
            CPU_MODEL=$(lscpu | grep "Model name" | cut -d':' -f2 | xargs || echo "Unknown")
            CPU_CORES=$(lscpu | grep "^CPU(s):" | cut -d':' -f2 | xargs || echo "Unknown")
            echo "  CPU: $CPU_MODEL ($CPU_CORES cores)"
        fi
        
        # Memory Info
        if [ -f /proc/meminfo ]; then
            TOTAL_RAM=$(grep MemTotal /proc/meminfo | awk '{printf "%.2f GB", $2/1024/1024}')
            echo "  Memory: $TOTAL_RAM"
        fi

        # GPU Info (NVIDIA check for Darknet)
        if command -v nvidia-smi >/dev/null 2>&1; then
            GPU_INFO=$(nvidia-smi --query-gpu=name --format=csv,noheader | head -n 1)
            echo "  GPU: $GPU_INFO (NVIDIA Driver active)"
        elif command -v lspci >/dev/null 2>&1; then
            if lspci | grep -qi "nvidia"; then
                echo "  GPU: NVIDIA hardware detected (driver might not be loaded)"
            fi
        fi

        # 7. Proactive Dependency Validation
        echo "[vms-lite] Validating installation and dependencies..."
        DEPS_OK=true
        
        # Check binary
        if [ -f "/usr/bin/vms-lite" ]; then
            echo "[OK] vms-lite binary found at /usr/bin/vms-lite"
        else
            echo "[ERROR] vms-lite binary NOT found at /usr/bin/vms-lite"
            DEPS_OK=false
        fi

        # Check bundled libraries and their dependencies
        for lib in "libdarkhelp.so" "libdarknet.so"; do
            lib_path="$VMS_LIBDIR/$lib"
            if [ -f "$lib_path" ]; then
                echo "[OK] Found $lib in private directory."
                
                # Check for missing system dependencies
                MISSING_DEPS=$(LD_LIBRARY_PATH="$VMS_LIBDIR" ldd "$lib_path" | grep "not found" || true)
                if [ -n "$MISSING_DEPS" ]; then
                    echo "[WARNING] $lib has missing dependencies on this system:"
                    echo "$MISSING_DEPS" | sed 's/^/    /'
                    echo "Please ensure all required system libraries (OpenCV, OpenBLAS, etc.) are installed."
                    echo "Recommended: sudo apt install libopencv-videoio4.x libopenblas0 libgomp1"
                    DEPS_OK=false
                fi
            else
                echo "[WARNING] $lib not found in $VMS_LIBDIR."
                # Don't set DEPS_OK=false here as they might be system-installed
            fi
        done

        if [ "$DEPS_OK" = true ]; then
            echo "[SUCCESS] All critical system dependencies verified."
        else
            echo "[CAUTION] Some dependencies are missing or validation failed. The service might fail to start."
        fi

        # 8. Handle Systemd service
        echo "Reloading systemd daemon..."
        systemctl daemon-reload || true

        if [ -z "$2" ]; then
            # New installation (not an upgrade)
            echo "Enabling vms-lite service..."
            systemctl enable vms-lite || true
            echo "Starting vms-lite service..."
            systemctl start vms-lite || true
        else
            # Upgrade
            echo "Restarting vms-lite service..."
            systemctl restart vms-lite || true
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
